# packaging/systemd/ports-agent.service
[Unit]
Description=Ports agent (port-forward / firewall control node)
Documentation=https://github.com/TheiVy-Service/ports-agent
After=network-online.target nftables.service
Wants=network-online.target

[Service]
Type=simple
ExecStart=/usr/bin/ports-agent
Restart=on-failure
RestartSec=5
User=root
ConfigurationDirectory=ports
StateDirectory=ports
LogsDirectory=ports
ReadWritePaths=/etc/ports /var/lib/ports /var/log/ports
AmbientCapabilities=CAP_NET_ADMIN CAP_NET_RAW
CapabilityBoundingSet=CAP_NET_ADMIN CAP_NET_RAW
NoNewPrivileges=true
ProtectHome=true
ProtectSystem=full

[Install]
WantedBy=multi-user.target
